# deploy/systemd/bloom-tochka-reconcile.timer.example
[Unit]
Description=Run Bloom Club Tochka reconciliation every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30
Persistent=true

[Install]
WantedBy=timers.target
